# grocery-store.service
[Unit]
Description=Grocery E-Store Application
After=network.target

[Service]
Type=simple
User=administrator
WorkingDirectory=/home/administrator/sari2store
ExecStart=/home/administrator/sari2store/start.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
